# data/patchs/lte/lte.service
[Unit]
Description=LTE

[Service]
Type=oneshot
ExecStart=/usr/bin/bg96 start
ExecStop=/usr/bin/bg96 stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target